# service_configs/math-backend.service
[Unit]
Description=Math Backend Service (FastAPI)
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/math
ExecStart=/home/ubuntu/math/backend/.venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
